# deploy/phantom.service
[Unit]
Description=Phantom Proxy Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/phantom-server /etc/phantom/server.toml
Restart=always
RestartSec=5
LimitNOFILE=65535

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/phantom

[Install]
WantedBy=multi-user.target
